# uwarl_summit_specific/60-can.rules
KERNEL=="pcanusb*", SYMLINK+="pcan_base", MODE="0666"
#RUN+="/bin/reset_can_non_xenomai.sh"
